# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
; Give php-fpm a moment to start before nginx begins accepting connections
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:prewarm]
; Warms the opcache after nginx and php-fpm are up by issuing a real PHP
; request. The response code doesn't matter — any response means PHP executed
; and opcache is loaded. Runs once at startup, never restarts.
command=/bin/sh -c "until curl -sf http://localhost/healthz > /dev/null 2>&1; do sleep 1; done && curl -s http://localhost/ > /dev/null 2>&1; true"
autostart=true
autorestart=false
startretries=0
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0